# com.punky.trading-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.punky.trading-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/manvinderarora/Desktop/trading_agent/crypto_trading_agent/run_forever.sh</string>
        <string>--live</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/manvinderarora/Desktop/trading_agent/crypto_trading_agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/manvinderarora/Desktop/trading_agent/crypto_trading_agent/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/manvinderarora/Desktop/trading_agent/crypto_trading_agent/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
