# launchd/com.manvinder.crypto-trading-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.manvinder.crypto-trading-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/manvinderarora/Library/Application Support/crypto_trading_agent_runtime/run_dashboard.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/manvinderarora/Library/Application Support/crypto_trading_agent_runtime</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>TRADING_DASHBOARD_PORT</key>
    <string>8080</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/manvinderarora/Library/Application Support/crypto_trading_agent_runtime/logs/dashboard.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/manvinderarora/Library/Application Support/crypto_trading_agent_runtime/logs/dashboard.stderr.log</string>
</dict>
</plist>
